# packaging/media-pi-agent.service
[Unit]
Description=Media Pi Agent REST Service
Documentation=https://github.com/sw-consulting/media-pi.device
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/media-pi-agent
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Allow systemctl reload to signal the main process to reload configuration.
# This uses HUP which the agent handles to reload `/etc/media-pi-agent/agent.yaml`.
ExecReload=/bin/kill -HUP $MAINPID

# Optionally, admins can provide the server key to ExecReload calls via
# an EnvironmentFile drop-in. Avoid placing secrets directly in the unit file.
# Example drop-in: /etc/systemd/system/media-pi-agent.service.d/override.conf
# [Service]
# EnvironmentFile=/etc/media-pi-agent/env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# Allow writes to the agent config dir, audio config and systemd unit/timer files
ReadWritePaths=/etc/media-pi-agent /etc/asound.conf /etc/systemd/system
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=media-pi-agent

[Install]
WantedBy=multi-user.target
